# packaging/systemd/neonfs-omnibus.service
[Unit]
Description=NeonFS Omnibus Daemon (Core + FUSE + NFS)
After=network.target
Documentation=https://neonfs.dev/docs

[Service]
Type=notify
User=neonfs
Group=neonfs
ExecStart=/usr/bin/neonfs-omnibus-daemon start
ExecStop=/usr/bin/neonfs-omnibus-daemon stop
Restart=on-failure
RestartSec=5
TimeoutStartSec=60
TimeoutStopSec=30

EnvironmentFile=-/etc/neonfs/neonfs.conf

RuntimeDirectory=neonfs
StateDirectory=neonfs
ConfigurationDirectory=neonfs

# FUSE requires CAP_SYS_ADMIN for mounting
AmbientCapabilities=CAP_SYS_ADMIN
NoNewPrivileges=false
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/neonfs /mnt

# This unit creates host-visible filesystem mounts (FUSE + NFS). The
# default mount-namespace propagation under PrivateTmp/ProtectSystem/
# ProtectHome is slave, so mounts created inside the unit don't escape
# back to the host namespace. `shared` makes the unit's mounts
# propagate outwards, at the cost of also seeing host-side mount
# changes — fine for a filesystem-providing service. (#758)
MountFlags=shared

[Install]
WantedBy=multi-user.target
