# packaging/systemd/neonfs.target
[Unit]
Description=NeonFS Distributed Filesystem
Documentation=https://neonfs.dev/docs
After=network.target
Requires=neonfs-core.service
Wants=neonfs-fuse.service neonfs-nfs.service

[Install]
WantedBy=multi-user.target
